# src/etc/ec2net/service/log/run
#!/bin/sh
exec chpst -u log svlogd /var/log/ec2net
